# ops/systemd/sitelayer-offregion-backup.service
[Unit]
Description=Sitelayer off-region prod backup (pg_dump -> DO Spaces)
Documentation=https://github.com/GitSteveLozano/sitelayer/blob/main/BACKUP_STRATEGY.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=sitelayer
Group=sitelayer
WorkingDirectory=/app/sitelayer
EnvironmentFile=/app/sitelayer/.env
ExecStart=/app/sitelayer/scripts/backup-to-offregion.sh
StandardOutput=journal
StandardError=journal
# Hardening: this unit only reads from PG, streams gzip, and writes to S3.
# It does not need write access to the rest of the filesystem.
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=/app/sitelayer
# pg_dump on a slow / busy cluster can take a while; cap at 4h so a hung
# session doesn't keep the timer blocked.
TimeoutStartSec=4h

[Install]
WantedBy=multi-user.target
