# ops/systemd/sitelayer-auto-deploy.service
[Unit]
Description=Sitelayer fleet-side auto-deploy watcher (dev+demo on origin/dev advance, never prod)
Documentation=https://github.com/GitSteveLozano/sitelayer/blob/main/docs/AUTO_DEPLOY.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Runs as the operator's user (user-level unit), so it reuses the existing ssh
# key + ssh config to the preview droplet (sitelayer@159.203.53.218) and the
# operator's git credentials. No doctl needed; just bash + curl + jq + git + ssh.
#
# WorkingDirectory is the operator's checkout for convenience only — the watcher
# itself never deploys from here; it maintains its OWN dedicated checkout under
# ~/.cache/sitelayer-autodeploy/repo (AUTODEPLOY_REPO_DIR).
WorkingDirectory=%h/projects/sitelayer
# PATH note: a user unit inherits a minimal PATH. Ensure git/curl/jq/ssh/flock
# resolve by including the usual operator bin dirs; adjust if your tools live
# elsewhere.
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/projects/sitelayer/scripts/fleet-auto-deploy.sh
StandardOutput=journal
StandardError=journal
# A single dev/demo watch-mode redeploy is quick, but a cold clone or a slow
# droplet can stretch; cap so a hung run can't wedge the timer indefinitely.
TimeoutStartSec=20min

[Install]
WantedBy=default.target
