# ops/systemd/sitelayer-e2e-runner.service
[Unit]
Description=Sitelayer deterministic e2e runner (verify:full on origin/dev + origin/main)
Documentation=https://github.com/GitSteveLozano/sitelayer/blob/main/docs/E2E_RUNNER.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Runs as the operator's user (user-level unit). Needs node/npm + docker on PATH
# because verify:full stands up the docker-compose app stack + Playwright. This
# is meant for a QUIET, idle box (the preview droplet off-hours or a $6/mo
# throwaway) — NOT taylor-pc, where a loaded machine makes e2e flake. See
# docs/E2E_RUNNER.md.
WorkingDirectory=%h/projects/sitelayer
# A user unit inherits a minimal PATH. Include the usual operator bin dirs plus
# wherever node/npm/docker resolve. Adjust if your toolchain lives elsewhere
# (e.g. an nvm shim dir, /usr/local/bin for docker).
Environment=PATH=%h/.local/bin:%h/.nvm/versions/node/current/bin:/usr/local/bin:/usr/bin:/bin
# VERIFY_LEVEL=full -> npm run verify:full (static+build+unit+integration+e2e).
# The runner sets this itself; left here as documentation of intent.
Environment=E2E_RUNNER_VERIFY_LEVEL=full
ExecStart=%h/projects/sitelayer/scripts/e2e-runner.sh
StandardOutput=journal
StandardError=journal
# verify:full pulls images, builds, and drives a browser — give it room but cap
# so a hung run (stuck browser/page-closed) can't wedge the timer indefinitely.
TimeoutStartSec=60min

[Install]
WantedBy=default.target
